# debian/postinst
#!/bin/sh

set -e

NAME=echessd

case "$1" in
    configure)
        chown root:$NAME /etc/echessd.conf
        chmod 0640 /etc/echessd.conf
        chown -R $NAME:$NAME /var/log/$NAME /var/lib/$NAME
        update-rc.d $NAME defaults
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
